# templates/etc/profile.d/vagrant-libvirt.sh.j2
# {{ ansible_managed }}
export VAGRANT_LIBVIRT_VIRT_SYSPREP_OPERATIONS="defaults,-ssh-userdir,-customize"
